# infra/cron/crontab
# In-repo cron table (#487).
#
# All entries call a wrapper script under /usr/local/bin/ rather than
# inlining curl. Wrappers own the env-file source, response-body
# capture, and `overall != 'ok'` detection — see #487 §3 for the
# contract.
#
# Cron interprets `%` as a newline literal, so any wrapper that needs a
# date format (e.g. for log filenames) MUST do its own
# `date +%Y%m%d-%H` *inside* the script, not in the crontab line.
#
# Triage baseline cadence — 15-minute per-customer fan-out
*/15 * * * * /usr/local/bin/run-triage-baseline-dispatch.sh

# Triage exclusion fanout — minute-scale drain of the durable queue
# (#457 ships the route, 1B-7 wires the cadence). Matches the worker's
# exponential-backoff schedule (1m / 5m / 25m / 2h / 12h) so the first
# retry tier fires on the next tick after a failure.
* * * * * /usr/local/bin/run-triage-exclusion-fanout.sh

# Corpus A retention — once daily at 03:15 UTC. Offset from the
# :00/:15/:30/:45 cadence ticks so the nightly run does not contend
# with a fan-out tick, and so an operator inspecting log_dir in the
# morning sees a recent run.
15 3 * * * /usr/local/bin/run-triage-baseline-retention.sh

# Corpus B retention + zombie reaper — every 6 hours. More frequent
# than corpus A because the zombie-runner reaper unblocks the active-
# fingerprint partial unique index for new runs of the same
# fingerprint; waiting a full day would leave a crashed runner's slot
# occupied for too long.
30 */6 * * * /usr/local/bin/run-triage-policy-retention.sh

# Triage condition-snapshot retention (#472) — once daily at 04:15 UTC,
# one hour after the corpus A sweep. The snapshot sweep prunes rows
# whose fingerprints are no longer referenced by either corpus, so the
# upstream corpus retention pass must complete first to release the
# references it controls. baseline_version_snapshot is retained
# forever and is skipped by this sweep.
15 4 * * * /usr/local/bin/run-triage-snapshot-retention.sh

# Triage engagement-signal retention (#588) — once daily at 05:15 UTC.
# Offset another hour after the snapshot sweep so all four retention
# jobs (baseline, policy, snapshot, engagement) write their logs in
# the same morning window without overlapping. The sweep is `created_at`-
# bounded and does not depend on the corpus / snapshot retention
# passes — the ordering is purely about log clustering.
15 5 * * * /usr/local/bin/run-triage-engagement-retention.sh
